# cookbooks/r2app/templates/default/supervisor.conf.erb
<% @programs.each do |name,spec| %>
[program:<%= name %>]
command=<%= spec[:command] %>
autorestart=<%= spec[:autorestart] %>
environment = TWISTD_ARGS="-n"
stdout_logfile=%(here)s/logs/%(program_name)s_stdout.log
stderr_logfile=%(here)s/logs/%(program_name)s_stderr.log

<% end %>

[unix_http_server]
file=%(here)s/supervisor.sock

[supervisord]
logfile=%(here)s/logs/supervisord.log
pidfile=%(here)s/supervisord.pid

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix://%(here)s/supervisor.sock
