# cookbooks/supervisord_shutdown/files/default/kill-supd.conf
kill timeout 70

start on runlevel RUNLEVEL=[016] PREVLEVEL=[2345]

task

exec /usr/local/bin/kill-supd.sh
