# cookbooks/epu/templates/default/supervisor.conf.erb
[program:<%= @epuservice[:program] %>]
command=<%= @epuservice[:command] %>
autorestart=<%= @epuservice[:autorestart] %>
stdout_logfile=%(here)s/logs/%(program_name)s_stdout.log
stderr_logfile=%(here)s/logs/%(program_name)s_stderr.log

[unix_http_server]
file=%(here)s/supervisor.sock

[supervisord]
logfile=%(here)s/logs/supervisord.log
pidfile=%(here)s/supervisord.pid

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix://%(here)s/supervisor.sock
